# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --save "" --appendonly no
autorestart=true
priority=10
stdout_logfile=/var/log/redis.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/redis.err
stderr_logfile_maxbytes=10MB

[program:migrate]
command=/app/.venv/bin/python manage.py migrate --noinput
directory=/app
autorestart=false
startsecs=0
exitcodes=0
priority=20
stdout_logfile=/var/log/migrate.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/migrate.err
stderr_logfile_maxbytes=10MB

[program:dramatiq]
command=/app/.venv/bin/python manage.py rundramatiq --processes 1 --threads 1
directory=/app
autorestart=true
stopwaitsecs=30
priority=30
stdout_logfile=/var/log/dramatiq.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/dramatiq.err
stderr_logfile_maxbytes=50MB

[program:gunicorn]
command=/app/.venv/bin/gunicorn --bind 0.0.0.0:8008 --workers 1 --threads 2 --worker-class gthread --timeout 600 musictranscription.wsgi
directory=/app
autorestart=true
priority=40
stdout_logfile=/var/log/gunicorn.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/gunicorn.err
stderr_logfile_maxbytes=50MB
